# root/etc/s6-overlay/s6-rc.d/init-mod-universal-tailscale/run
#!/usr/bin/with-contenv bash

## Ubuntu
if [ -f /usr/bin/apt ]; then
    echo "\
        iptables" >> /mod-repo-packages-to-install.list

fi
# Alpine
if [ -f /sbin/apk ]; then
    echo "\
        iptables" >> /mod-repo-packages-to-install.list
fi
